# files/etc/uci-defaults/90-pcat2-setup
#!/bin/sh

# Enable and start the display service
if [ -x /etc/init.d/pcat2-display-mini ]; then
    /etc/init.d/pcat2-display-mini enable
    /etc/init.d/pcat2-display-mini start
fi

exit 0
